# systemd/mesh-monitor-collector.service
[Unit]
Description=Mesh Network Monitor - Metrics Collector
After=network.target mesh-monitor.service

[Service]
Type=simple
User=mesh-monitor
Group=mesh-monitor
WorkingDirectory=/opt/mesh-monitor
ExecStart=/opt/mesh-monitor/venv/bin/python3 /opt/mesh-monitor/collector.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
