# systemd/unbound-cache-sync.service
[Unit]
Description=Unbound DNS Cache Sync via etcd
After=network.target etcd.service unbound.service
Requires=etcd.service unbound.service

[Service]
Type=simple
ExecStart=/usr/local/bin/unbound-cache-sync.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
